# templates/bsl_bootstrap.debian.init.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          bsl_buppet
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: bsl_bootstrap provisioning
# Description:       Startup Puppet provisioning using bsl_puppet::setup
### END INIT INFO

. /lib/lsb/init-functions

set +e

LOGFILE='<%= scope.lookupvar('bsl_bootstrap::puppetmaster::config::init_service_puppet_log') %>'
PUPPET='<%= scope.lookupvar('bsl_bootstrap::puppet_binary') %>'

do_start () {
  log_begin_msg "bsl_bootstrap startup provisioning begins via ${PUPPET}"

  <%= scope.lookupvar('bsl_bootstrap::puppetmaster::config::init_service_facter_vars') -%> $PUPPET apply \
    <%= scope.lookupvar('bsl_bootstrap::puppetmaster::config::init_service_puppet_args') %> \
    --certname='<%= scope.lookupvar('bsl_bootstrap::puppetmaster::config::target_certname') %>' \
    --environment='<%= scope.lookupvar('bsl_bootstrap::puppetmaster::config::environment') %>' \
    --logdest=$LOGFILE \
    /etc/puppetlabs/code/environments/<%= @environment %>/manifests/site.pp

  log_end_msg $?
}

case "$1" in
  start)
    do_start
    exit $?
  ;;
  *)
  echo "Usage: $0 {start}" >&2
  exit 1
  ;;
esac
